# ytuner/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: YTuner - Nginx Reverse Proxy
# ==============================================================================

bashio::log.info "Starting nginx reverse proxy..."
exec nginx -g 'daemon off;'
