# ytuner/rootfs/etc/s6-overlay/s6-rc.d/svc-ytuner/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: YTuner - YTuner Server
# Logs to /data/ytuner.log so the web UI discovery can tail it
# ==============================================================================

YTUNER_LOG="/data/ytuner.log"

bashio::log.info "Starting YTuner server (log: ${YTUNER_LOG})..."
cd /opt/ytuner || exit 1
/opt/ytuner/ytuner >> "${YTUNER_LOG}" 2>&1
